# tests/fixtures/source/etc/systemd/system/argo.service
[Unit]
Description=Cloudflare Tunnel

[Service]
ExecStart=/etc/sing-box/argo tunnel --url http://localhost:8001
