# tests/fixtures/source/etc/systemd/system/sing-box.service
[Unit]
Description=sing-box service

[Service]
ExecStart=/etc/sing-box/sing-box run -C /etc/sing-box/conf

